# board/teknotelelectronics/maestro-cpx/overlay/var/spool/cron/crontabs/root
SHELL=/bin/bash
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

#Set the hardware clock
00 03 * * * /sbin/hwclock --systohc

#Maintain user services and apps
* * * * * /usr/local/cpx/scripts/maintain.sh
