# board/teknotelelectronics/maestro-cpx/overlay/var/spool/cron/crontabs/root
SHELL=/bin/bash
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

#Maintain user services and apps every minute
* * * * * /usr/local/cpx/scripts/maintain.sh maintain

#Synchronize the hardware clock at 03:00
00 03 * * * /usr/local/cpx/scripts/maintain.sh syncTime